# service/usernetes-control-plane.service
[Unit]
Description=User Space Kubernetes Control Plane (Usernetes)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/workspace/usernetes/usernetes-ipa/service/usernetes-start-control-plane.sh podman
LimitMEMLOCK=infinity
# Environment="USERNETES_CONTAINER_TECH_ARG=podman"

# Standard output and error
# Customize this to where you want it.  
StandardOutput=file:/usr/workspace/usernetes/control-plane.log
StandardError=file:/usr/workspace/usernetes/control-plane.log

# KillMode=process is important. When the service stops, systemd will only kill
# the main process (sleep infinity). The Kubernetes components run inside containers
# managed by Podman/Docker, which `make down-v` (or manual podman commands) should handle
# on a clean shutdown. If the node is just "cleaned up" (deallocated), then everything dies anyway.
KillMode=process

# Start when the user's session starts
[Install]
WantedBy=default.target
